# deploy/systemd/lifelog-server.service
# Dev/build-from-source service template
[Unit]
Description=Lifelog Server Backend
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/Projects/lifelog
Environment=HOME=/home/YOUR_USER
Environment=LIFELOG_HOST=0.0.0.0
# DB credentials loaded from EnvironmentFile — do not hardcode here
Environment=LIFELOG_POSTGRES_INGEST_URL=postgresql://lifelog@127.0.0.1:5432/lifelog
Environment=LIFELOG_POSTGRES_INGEST_MAX_CONNECTIONS=16
Environment=LIFELOG_CONFIG_PATH=/home/YOUR_USER/Projects/lifelog/lifelog-config.toml
Environment=RUST_LOG=info
ExecStart=/usr/bin/env bash -lc 'exec nix develop --command cargo run -p lifelog-server --bin lifelog-server'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
