# deploy/systemd/lifelog-surrealdb.service
[Unit]
Description=Lifelog SurrealDB
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=matth
WorkingDirectory=/home/matth/Projects/lifelog
Environment=HOME=/home/matth
ExecStartPre=/usr/bin/mkdir -p /home/matth/lifelog/surrealdb
ExecStart=/usr/bin/env bash -lc 'exec nix develop --command surreal start --bind 127.0.0.1:7183 --user root --pass root file:/home/matth/lifelog/surrealdb/main.db'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
