# deploy/systemd-user/lifelog-surrealdb.service
[Unit]
Description=Lifelog SurrealDB (user)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Projects/lifelog
Environment=HOME=%h
ExecStartPre=/usr/bin/env bash -lc 'mkdir -p %h/lifelog/surrealdb'
ExecStart=/usr/bin/env bash -lc 'exec nix develop --command surreal start --bind 127.0.0.1:7183 --user root --pass root file:%h/lifelog/surrealdb/main.db'
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
